# 12272020/root/opt/arklone/systemd/arkloned-roms.path
[Path]
PathChanged=/roms
Unit=arkloned@retroarch-roms\x40-roms.service
