# 12272020/root/opt/arklone/systemd/arkloned-retroarch-contentroot.path
[Path]
PathChanged=/roms
Unit=arkloned@-roms\x40retroarch-roms\x40retroarch.service

[Install]
WantedBy=multi-user.target
