# 12272020/root/opt/arklone/systemd/arkloned-states.path
[Path]
PathChanged=/home/ark/.config/retroarch/states
Unit=arkloned@retroarch-states\x40-home-ark-.config-retroarch-states.service
