# 12272020/root/opt/arklone/systemd/arkloned-amiberry.path
[Path]
PathChanged=/opt/amiberry/savestates
Unit=arkloned@amiberry-savestates\x40-opt-amiberry-savestates.service
