# 12272020/root/opt/arklone/systemd/arkloned-VVVVVV.path
[Path]
PathChanged=/roms/ports/VVVVVV/saves
Unit=-roms-ports-VVVVVV-saves\x40VVVVVV-saves\x40VVVVVV.service

[Install]
WantedBy=multi-user.target
Wants=network-online.target
After=network-online.target
