# 12272020/root/opt/arklone/systemd/arkloned-saves.path
[Path]
PathChanged=/home/ark/.config/retroarch/saves
Unit=arkloned@retroarch-saves\x40-home-ark-.config-retroarch-saves.service
